# zookeeper/supervisor.ini
[program:zookeeper]
command=/srv/zookeeper/bin/zkServer.sh start-foreground
redirect_stderr=false
stdout_logfile=/var/log/zookeeper/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/zookeeper/stderr
stderr_logfile_maxbytes=0
stopsignal=INT
numprocs_start=1
startsecs=2
autostart=true
autorestart=true
