# kafka/supervisor-connect.ini
[program:connect]
command=/opt/kafka/bin/connect-distributed.sh /opt/kafka/config/connect-distributed.properties
redirect_stderr=false
stdout_logfile=/var/log/connect/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/connect/stderr
stderr_logfile_maxbytes=0
stopsignal=INT
numprocs_start=1
startsecs=2
autostart=true
autorestart=true
