# kafka/supervisor.ini
[program:kafka]
command=/opt/kafka/bin/kafka-server-start.sh /opt/kafka/config/server.properties --override zookeeper.connect="%(ENV_ZOOKEEPER_HOSTS)s"
redirect_stderr=false
stdout_logfile=/var/log/kafka/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/kafka/stderr
stderr_logfile_maxbytes=0
stopsignal=INT
numprocs_start=1
startsecs=2
autostart=true
autorestart=true
